# supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log   ; supervisord log file
pidfile=/var/run/supervisord.pid              ; pidfile location
childlogdir=/var/log/supervisor/              ; where child log files will live
logfile_maxbytes=50MB                         ; maximum size of logfile before rotation
logfile_backups=10                            ; number of backed up logfiles
loglevel=info                                 ; info, debug, warn, trace
minfds=1024                                   ; number of startup file descriptors
minprocs=200                                  ; number of process descriptors
user=root                                     ; default user

[unix_http_server]
file=/tmp/supervisor.sock   ; path to your socket file

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use unix:// schem for a unix sockets.

[program:nginx]
command = /usr/sbin/nginx 
user = root
autostart = true

[program:postgresql]
command=/usr/lib/postgresql/9.6/bin/postgres -D /var/lib/postgresql/9.6/main -c config_file=/etc/postgresql/9.6/main/postgresql.conf
#user=postgres
user=%(ENV_APPUSER)s
autostart=true
autorestart=true

[program:mongodb]
command=/usr/bin/mongod --dbpath /var/lib/mongodb
;directory=/var/lib/mongodb
user=root
autostart=true

[program:sshd]
command=/usr/sbin/sshd -D
numprocs=1
autostart=true
autorestart=true

[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server $@
user=%(ENV_APPUSER)s
numprocs=1
autostart=true
autorestart=true

[include]
files=/etc/supervisor/conf.d/*.conf 
